# deploy/celery-worker.service
[Unit]
Description=Trading Intelligence Celery Worker
After=network.target redis-server.service
PartOf=trading-stack.target

[Service]
Type=simple
User=trading
Group=trading
WorkingDirectory=/app/profitlyworkflow
EnvironmentFile=/app/profitlyworkflow/.env.production
ExecStart=/app/profitlyworkflow/venv/bin/celery -A flows.celery_app worker --loglevel=info --concurrency=2
Restart=always
RestartSec=10

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/app/profitlyworkflow

[Install]
WantedBy=multi-user.target
